# src/packages/fff/fff-fastd/files/etc/uci-defaults/55_fff-fastd
/etc/init.d/fastd disable

project="$VPN_PROJECT"

>/etc/config/fastd

uci batch <<EOF
  set fastd.FFF='fastd'
  set fastd.FFF.enabled='1'
  set fastd.FFF.config_peer_dir="/etc/fastd/FFF/peers"
  set fastd.FFF.syslog_level='warn'
  set fastd.FFF.method='null'
  set fastd.FFF.mode='tap'
  set fastd.FFF.interface="FFFVPN"
  set fastd.FFF.mtu='1426'
  set fastd.FFF.on_up="/etc/fastd/FFF/up.sh"
  set fastd.FFF.secure_handshakes='0'
  set fastd.FFF.secret="generate"
EOF

[ ! -d /etc/fastd/FFF ] &&  mkdir -p /etc/fastd/FFF
ln -s /tmp/fastd_FFF_peers /etc/fastd/FFF/peers
echo "#!/bin/sh" > /etc/fastd/FFF/up.sh
echo "ip link set up dev FFFVPN" >> /etc/fastd/FFF/up.sh
echo "echo enable > /sys/devices/virtual/net/FFFVPN/batman_adv/no_rebroadcast" >> /etc/fastd/FFF/up.sh
echo "batctl if add FFFVPN" >> /etc/fastd/FFF/up.sh
chmod +x /etc/fastd/FFF/up.sh
